# sshmux.service
[Unit]
Description=sshmux SSH multiplexer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sshmux
Group=sshmux
WorkingDirectory=/var/lib/sshmux
ExecStart=/usr/local/bin/sshmux --config /etc/sshmux/config.yaml 0.0.0.0:22
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sshmux

[Install]
WantedBy=multi-user.target
